# deploy/supervisor-wsl.conf
; Velour supervisor config for local WSL development.
;
; Manages the dev server and periodic jobs so Velour is "always on"
; whenever a WSL session is open. Install:
;
;   sudo apt-get install -y supervisor
;   sudo cp deploy/supervisor-wsl.conf /etc/supervisor/conf.d/velour-wsl.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Or start a user-mode supervisord (no sudo needed):
;
;   supervisord -c deploy/supervisord-wsl.ini
;
; Logs go to /tmp/velour_*.log so they survive across sessions but
; don't accumulate silently in /var/log.

[program:velour-runserver]
command=/home/handyc/claubsh/velour-dev/venv/bin/python
    /home/handyc/claubsh/velour-dev/manage.py smart_runserver
directory=/home/handyc/claubsh/velour-dev
user=handyc
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/velour_runserver.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
environment=DJANGO_SETTINGS_MODULE="velour.settings"
startsecs=3
stopwaitsecs=5

[program:velour-identity-cron]
; Fires every 60 seconds via a shell loop. Replaces the crontab entry
; so everything is under one process manager.
command=/bin/bash -c "while true; do /home/handyc/claubsh/velour-dev/venv/bin/python /home/handyc/claubsh/velour-dev/manage.py identity_cron --quiet; sleep 60; done"
directory=/home/handyc/claubsh/velour-dev
user=handyc
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/velour_identity_cron.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
environment=DJANGO_SETTINGS_MODULE="velour.settings"
startsecs=2

[program:velour-poll-inbox]
; Poll inbound mail every 5 minutes.
command=/bin/bash -c "while true; do /home/handyc/claubsh/velour-dev/venv/bin/python /home/handyc/claubsh/velour-dev/manage.py poll_inbox 2>&1 || true; sleep 300; done"
directory=/home/handyc/claubsh/velour-dev
user=handyc
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/velour_poll_inbox.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
environment=DJANGO_SETTINGS_MODULE="velour.settings"
startsecs=2
